# OpenWrt/etc/init.d/mpu
#!/bin/sh /etc/rc.common
# Copyright (C) 2008 OpenWrt.org

START=97

start() {
sleep 10
mpu6050 > /tmp/mpu6050.txt &
sleep 10
/bin/sh /bin/6050.sh > /www/mess.log &
}